# claude-certified-architect-lab.service
[Unit]
Description=Claude Certified Architect Lab API
After=network.target

[Service]
Type=simple
User=claudeapp
Group=claudeapp
WorkingDirectory=/opt/claude-certified-architect-lab
EnvironmentFile=/etc/claude-certified-architect-lab/claude-certified-architect-lab.env
Environment=PYTHONPATH=/opt/claude-certified-architect-lab
ExecStart=/opt/claude-certified-architect-lab/.venv/bin/python -m uvicorn src.main:app \
    --host 0.0.0.0 \
    --port 8000 \
    --workers 4
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-certified-architect-lab

[Install]
WantedBy=multi-user.target
